# systemd/hmdistreamer-hdmi-bringup.service
[Unit]
Description=Configure X1300 HDMI capture pipeline
Before=hmdistreamer-ndi-sender.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
EnvironmentFile=-/etc/hmdistreamer/hmdistreamer.env
ExecStart=/usr/local/bin/hmdistreamer-hdmi-bringup
RemainAfterExit=yes
Restart=on-failure
RestartSec=2
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
